# .n3rv/systemd/nerv-hub.service
[Unit]
Description=NERV A2A Hub for l06_p0s3
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ultravietnamita/Projects/roofs/l06_p0s3
ExecStart=/home/ultravietnamita/Projects/nerv/.venv/bin/nerv hub start
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ultravietnamita/Projects/roofs/l06_p0s3/.nerv/logs/hub.log
StandardError=append:/home/ultravietnamita/Projects/roofs/l06_p0s3/.nerv/logs/hub.log

[Install]
WantedBy=default.target