# .n3rverberage/systemd/n3rverberage-hub.service
[Unit]
Description=N3RVERBERAGE A2A Hub for hub
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/ultravietnamita/Projects/reverberage/hub
ExecStart=/home/ultravietnamita/.local/bin/n3rverberage hub start
Restart=on-failure
RestartSec=5
StandardOutput=append:/home/ultravietnamita/Projects/reverberage/hub/.n3rverberage/logs/hub.log
StandardError=append:/home/ultravietnamita/Projects/reverberage/hub/.n3rverberage/logs/hub.log

[Install]
WantedBy=default.target